# stubs/crontab
# CronTab for running laravels scheduler
* * * * * php /project/artisan schedule:run >> /var/log/cron.log
